# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
autorestart=true
command=/opt/start_app.sh

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true